# services.d/nginx/run
#!/usr/bin/with-contenv sh

j2 /etc/nginx/imported/default.conf.j2 > /etc/nginx/conf.d/default.conf

exec nginx
